# root/etc/s6-overlay/s6-rc.d/10-create-exports/finish
#!/command/with-contenv /bin/sh

/usr/sbin/exportfs -uav
/run/s6/basedir/bin/halt
