# root/etc/s6-overlay/s6-rc.d/40-exportfs-log/run
#!/command/with-contenv /bin/sh
program_name="exportfs"

s6-format-filter "%1 %s" "[${program_name}]" | logutil-service /log/"${program_name}"
