# root/etc/s6-overlay/s6-rc.d/mountd/finish
#!/command/with-contenv sh

/run/s6/basedir/bin/halt
