# root/etc/s6-overlay/s6-rc.d/20-rpcbind-log/run
#!/command/with-contenv sh

program_name="rpcbind"

s6-format-filter "%1 %s" "[${program_name}]" | s6-log -b -- n10 s1000000 T 1 T /log/"${program_name}"
