# root/etc/s6-overlay/s6-rc.d/20-rpcbind/finish
#!/command/with-contenv /bin/sh

/run/s6/basedir/bin/halt
