# root/etc/s6-overlay/s6-rc.d/30-nfsd/finish
#!/command/with-contenv sh

/usr/sbin/rpc.nfsd 0
/run/s6/basedir/bin/halt
